# contrib/docker/fabmanager.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/rediserror.log

[program:postgresql]
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
autostart=true
autorestart=true
user=postgres
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgreserror.log

[program:fabmanager]
command=su -l fabmanager -c "cd /home/fabmanager/fab-manager && foreman start"
directory=/home/fabmanager/fab-manager
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/fabmanager.log
stderr_logfile=/var/log/supervisor/fabmanagererror.log